# psplash/initrd/scripts/init-bottom/end.psplash.service
[Unit]
Description=Stop Splash Screen

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill psplash

[Install]
WantedBy=multi-user.target
